# prometheus-daemon.service
[Unit]
Description=Prometheus Market-Aware Orchestration Daemon (new stack)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=feanor
Group=feanor
WorkingDirectory=/home/feanor/coding/prometheus
EnvironmentFile=/etc/sysconfig/prometheus-daemon

# Use the system Python binary (bin_t, SELinux-compatible) with the venv
# site-packages on PYTHONPATH so imports resolve without needing the symlinks
# inside ~/coding/prometheus/.venv/bin/ (which are user_home_t and not exec-able
# by systemd's init domain).
Environment=PYTHONPATH=/home/feanor/coding/prometheus/.venv/lib/python3.14/site-packages:/home/feanor/coding/apathis:/home/feanor/coding/prometheus

# Market-aware daemon: polls every 60s, executes full DAG for US_EQ / KRONOS / INTEL
ExecStart=/usr/bin/python3.14 \
    -m prometheus.orchestration.market_aware_daemon \
    --market US_EQ \
    --market KRONOS \
    --market INTEL \
    --options-mode paper \
    --poll-interval-seconds 60

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=90

# Restart on failure with backoff
Restart=on-failure
RestartSec=30

# Security
NoNewPrivileges=yes

# Logging (goes to systemd journal; read with: journalctl -u prometheus-daemon -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prometheus-daemon

[Install]
WantedBy=multi-user.target
